# external/board/beaglebone/overlay/usr/lib/systemd/system/mender-client.service
[Unit]
Description=Mender OTA update service
After=systemd-resolved.service resize_data_part.service

[Service]
Type=idle
User=root
Group=root
ExecStart=/usr/bin/mender -daemon
Restart=on-abort

[Install]
WantedBy=multi-user.target
