# external/board/beaglebone/overlay/usr/lib/systemd/system/resize_data_part.service
[Unit]
Description=resize data partition
Before=data.mount setup_directories.service

[Service]
Type=oneshot
ExecStart=/usr/bin/resize_data_part
KillSignal=SIGQUIT
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
